# vultr_setup/verzekbot.service
[Unit]
Description=Verzek Telegram Signal Forwarder
After=network.target

[Service]
Type=simple
ExecStart=/var/www/VerzekAutoTrader/venv/bin/python3 /var/www/VerzekAutoTrader/telethon_forwarder.py
WorkingDirectory=/var/www/VerzekAutoTrader
User=root
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
