# signal_engine/systemd/verzek-signalengine.service
[Unit]
Description=VerzekSignalEngine - Multi-Bot Trading Signal System
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/signal_engine
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/root/signal_engine"
ExecStart=/usr/bin/python3 -m main
Restart=always
RestartSec=10
StandardOutput=append:/root/signal_engine/logs/systemd.log
StandardError=append:/root/signal_engine/logs/systemd_error.log

# Resource limits
LimitNOFILE=65536

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
