# backend/deploy/verzek_worker.service
[Unit]
Description=Verzek AutoTrader Worker
After=network.target verzek_api.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/VerzekBackend/backend
EnvironmentFile=/root/VerzekBackend/backend/.env
ExecStart=/usr/bin/python3 worker.py
Restart=always
RestartSec=10
StandardOutput=append:/root/VerzekBackend/backend/logs/worker.log
StandardError=append:/root/VerzekBackend/backend/logs/worker_error.log

[Install]
WantedBy=multi-user.target
