# backend/deploy/verzek_worker.service
[Unit]
Description=Verzek AutoTrader Worker
After=network.target verzek_api.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/api_server/backend
EnvironmentFile=/root/api_server_env.sh
ExecStart=/root/api_server/backend/venv/bin/python worker.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
