# backend/deploy/verzek_api.service
[Unit]
Description=Verzek AutoTrader API Server
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/VerzekBackend/backend
EnvironmentFile=/root/api_server_env.sh
ExecStart=/usr/local/bin/gunicorn -c gunicorn.conf.py api_server:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
