# backend/deploy/verzek_api.service
[Unit]
Description=Verzek AutoTrader API Server
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/VerzekBackend/backend
EnvironmentFile=/root/VerzekBackend/backend/.env
ExecStart=/usr/local/bin/gunicorn -c gunicorn.conf.py api_server:app
Restart=always
RestartSec=10
StandardOutput=append:/root/VerzekBackend/backend/logs/api_server.log
StandardError=append:/root/VerzekBackend/backend/logs/api_server_error.log

[Install]
WantedBy=multi-user.target
